# config/syspilot.service
[Unit]
Description=System Health Monitor Daemon
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/path/to/syspilot/sys/main.py
WorkingDirectory=%h/path/to/syspilot/sys
StandardOutput=journal
StandardError=journal
Environment="PATH=/usr/bin:/usr/local/bin"
RemainAfterExit=no

[Install]
WantedBy=default.target
